# etc/supervisor/conf.d/foo.conf
; http://supervisord.org/running.html#adding-a-program
[program:foo]
command=/bin/cat
